# superset/superset.service
[Unit]
Description=Superset
[Service]
ExecStart=/home/ubuntu/superset.sh
[Install]
WantedBy=multi-user.target