# superset/superset.service
[Unit]
Description=Superset
[Service]
ExecStart=/usr/bin/superset.sh
[Install]
WantedBy=multi-user.target